# systemd/dftool.service
[Unit]
Description=DFTool - Digital Forensics Monitoring Daemon
Documentation=https://github.com/yourname/dftool
After=network.target syslog.target
Wants=network.target

[Service]
Type=forking
PIDFile=/var/run/dftool/dftoold.pid
ExecStart=/usr/local/bin/dftoold start
ExecStop=/usr/local/bin/dftoold stop
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dftool

# Security hardening (daemon itself runs as root to read all logs)
ProtectSystem=strict
ReadWritePaths=/var/log/dftool /var/lib/dftool /var/run/dftool
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
